# docker/app/etc/services.d/app/run
#!/usr/bin/with-contenv sh
cd /usr/src/app
exec npm start
